# recipes-efl/e17/files/enlightenment.service.in
[Unit]
Description=Enlightenment service
After=dbus.socket

[Service]
Environment=XDG_RUNTIME_DIR=/tmp
Environment=HOME=/home/root/
Environment=E_START=enlightenment
ExecStart=@prefix@/bin/enlightenment
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
